# rootfs/etc/sv/sysctl/run
#!/bin/sh
sleep 15
sysctl -p
exec chpst -b sysctl pause
