# rootfs/etc/cron.monthly/fstrim
#!/bin/sh
fstrim -Av | logger -t fstrim -p daemon.info
