# rootfs/etc/cron.monthly/unbound
#!/bin/sh
curl -fs -m 60 -o /var/lib/unbound/root.hints https://www.internic.net/domain/named.root
