# rootfs-zfs/etc/cron.monthly/fstrim
#!/bin/sh
fstrim -v /boot | logger -t fstrim -p daemon.info
